# systemd/csr-popup.service.in
[Unit]
Description=CSR popup service for user session

[Service]
User=app
Group=app
Type=simple
SmackProcessLabel=@SERVICE_NAME@
EnvironmentFile=/run/tizen-mobile-env
ExecStart=@BIN_DIR@/@SERVICE_NAME@-popup
Sockets=@SERVICE_NAME@-popup.socket

[Install]
WantedBy=default.target
