# systemd/csr-popup.service.in
[Unit]
Description=CSR popup service for user session

[Service]
@DECLARE_POPUP_USER@
@DECLARE_POPUP_GROUP@
@DECLARE_POPUP_SMACK_PROCESS_LABEL@
Type=simple
EnvironmentFile=@POPUP_SERVICE_ENV_FILE_PATH@
ExecStart=@BIN_DIR@/@SERVICE_NAME@-popup
Sockets=@SERVICE_NAME@-popup.socket
