# systemd/csr.service.in
[Unit]
Description=Content Screening and Reputation service
Requires=@SERVICE_NAME@.socket

[Service]
User=@SERVICE_USER@
Group=@SERVICE_GROUP@
SmackProcessLabel=@SERVICE_NAME@
ExecStart=@BIN_DIR@/@SERVICE_NAME@-server
Sockets=@SERVICE_NAME@.socket

[Install]
WantedBy=multi-user.target
