# systemd/csr.service.in
[Unit]
Description=Content Screening and Reputation service
Requires=@SERVICE_NAME@-cs.socket @SERVICE_NAME@-wp.socket @SERVICE_NAME@-admin.socket

[Service]
User=@SERVICE_USER@
Group=@SERVICE_GROUP@
SmackProcessLabel=@SERVICE_NAME@
ExecStart=@BIN_DIR@/@SERVICE_NAME@-server
Sockets=@SERVICE_NAME@-cs.socket
Sockets=@SERVICE_NAME@-wp.socket
Sockets=@SERVICE_NAME@-admin.socket

[Install]
WantedBy=multi-user.target
