# systemd/csr.socket.in
[Unit]
Description=Csr framework socket
Wants=@SERVICE_NAME@.service
Before=@SERVICE_NAME@.service

[Socket]
ListenStream=/tmp/.@SERVICE_NAME@.socket
Service=@SERVICE_NAME@.service
SocketMode=0777
# smack access control on socket disabled shortly to fast development
#SmackLabelIPIn=@SERVICE_NAME@::api
SmackLabelIPIn=*
SmackLabelIPOut=@

[Install]
WantedBy=sockets.target
